# stellar-trading-bot.service
[Unit]
Description=Stellar XLM/USDC Trading Bot
After=network.target postgresql.service stellar-data-collector.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/tradebot
ExecStart=/usr/bin/python3 /home/ec2-user/tradebot/src/main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
